# scripts/supervisor-expense-app.conf
; Supervisor configuration for Expense App
; Copy to: /etc/supervisor/conf.d/expense-app.conf

[program:expense-backend]
process_name=%(program_name)s
command=php /var/www/expense-app/backend/artisan serve --host=127.0.0.1 --port=8000
directory=/var/www/expense-app/backend
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/expense-backend.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600
stopsignal=SIGTERM

[program:expense-queue]
process_name=%(program_name)s
command=php /var/www/expense-app/backend/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/expense-app/backend
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/expense-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600
stopsignal=SIGTERM

[group:expense-app]
programs=expense-backend,expense-queue
priority=999
